# notion-bridge.service
[Unit]
Description=AI Dala Notion Bridge Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/chatgpt_bridge
EnvironmentFile=/root/chatgpt_bridge/.env
ExecStart=/root/chatgpt_bridge/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target